# ossxp-linux/debian/ossxp-linux-base.conffiles
/etc/cron.d/ossxp-backup
/opt/ossxp/conf/ssl/openssl.cnf
